# ansible/roles/gateway/templates/wg-quick.init.j2
#!/sbin/openrc-run
# OpenRC init script for wg-quick@wg0.
# Requires /etc/wireguard/wg0.conf — written by wireguard-setup.sh in cloud-init.

description="WireGuard VPN interface wg0"

command=/usr/bin/wg-quick
command_args="up wg0"
stop_command="/usr/bin/wg-quick down wg0"

depend() {
    need net
    after cloud-final
}

stop() {
    /usr/bin/wg-quick down wg0 2>/dev/null || true
}
